# kanata/.config/systemd/user/kanata.service
[Unit]
Description=kanata keyboard remapper
Documentation=https://github.com/jtroo/kanata
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/kanata --nodelay --cfg %h/.config/kanata/kanata.kbd
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
